# Downloads/claude-code-ssh-master/s6-overlay/s6-rc.d/fail2ban/run
#!/command/execlineb -P
fdmove -c 2 1
fail2ban-server -xf start
